# share/droneshot-daemon/droneshot-daemon.service.template
[Unit]
Description=Drone Shot System

[Install]
WantedBy=multi-user.target

[Service]
Type=notify
ExecStart=PATH_TO_EXE
WorkingDirectory=PATH_TO_DIST
FailureAction=reboot
KillSignal=SIGINT
