# share/droneshot-daemon/droneshot-daemon.service.template
[Service]
Type=notify
ExecStart=PATH_TO_EXE
WorkingDirectory=PATH_TO_DIST
FailureAction=reboot
KillSignal=SIGINT
